# docker/app/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; -------------------------------------------------------------------
; Глобальные настройки для всех программ: рабочая директория /app
; -------------------------------------------------------------------

[program:post-startup]
command=/usr/local/bin/post-startup.sh
autostart=true
autorestart=false          ; скрипт должен выполниться один раз
startsecs=5                ; увеличено с 0, чтобы избежать ложного успеха
stderr_logfile=/app/storage/logs/post-startup.err.log
stdout_logfile=/app/storage/logs/post-startup.out.log
user=root
priority=0
directory=/app             ; явно задаём рабочую директорию
stopwaitsecs=10

[program:octane]
command=/usr/local/bin/start_server.sh
autostart=true
autorestart=true
stdout_logfile=/app/storage/logs/octane.out.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/app/storage/logs/octane.err.log
stderr_logfile_maxbytes=50MB
user=root
priority=1
directory=/app
stopwaitsecs=10

[program:ssr]
command=php /app/artisan inertia:start-ssr
autostart=true
autorestart=true
stdout_logfile=/app/storage/logs/ssr.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/app/storage/logs/ssr.err.log
stderr_logfile_maxbytes=10MB
user=root
priority=1
directory=/app
stopwaitsecs=10
startsecs=3

; Шаблон для всех queue-воркеров:
; - run-worker-loop.sh должен принимать команду как аргументы (без лишних кавычек)
; - В самом скрипте желательно использовать exec "$@", чтобы supervisor управлял процессом.
; - Если скрипт всё же ожидает один строковый аргумент, тогда оставьте кавычки.
;   Приведённый ниже вариант предполагает корректную передачу аргументов.

[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/run-worker-loop.sh "php /app/artisan queue:work --sleep=3 --tries=3 --max-time=3600"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=5
redirect_stderr=true
stdout_logfile=/app/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
priority=2
directory=/app
stopwaitsecs=10

[program:schedule]
command=/usr/local/bin/run-worker-loop.sh "php /app/artisan schedule:work"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
user=root
priority=3
directory=/app
stopwaitsecs=10

; [program:reverb]
; command=/usr/local/bin/run-worker-loop.sh "php /app/artisan reverb:start --host=0.0.0.0 --port=8080"
; autostart=true
; autorestart=true
; redirect_stderr=true
; stdout_logfile=/app/storage/logs/reverb.log
; stdout_logfile_maxbytes=10MB
; stdout_logfile_backups=5
; user=root
; priority=5
; directory=/app
; stopwaitsecs=10

[program:prune]
command=/usr/local/bin/run-worker-loop.sh "php /app/artisan queue:work --queue=prune"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/queue-prune.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
priority=4
directory=/app
stopwaitsecs=10

[program:content]
process_name=content_%(program_name)s_%(process_num)02d
command=/usr/local/bin/run-worker-loop.sh "php /app/artisan queue:work --queue=content --sleep=3 --tries=3 --max-time=3600"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=2
redirect_stderr=true
stdout_logfile=/app/storage/logs/worker.content.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
priority=2
directory=/app
stopwaitsecs=10

[program:base]
command=/usr/local/bin/run-worker-loop.sh "php /app/artisan queue:work --queue=base"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/queue-base.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
priority=4
directory=/app
stopwaitsecs=10

[program:work1]
command=/usr/local/bin/run-worker-loop.sh "php /app/artisan queue:work --queue=work1"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/queue-work1.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
priority=4
directory=/app
stopwaitsecs=10

[program:work2]
command=/usr/local/bin/run-worker-loop.sh "php /app/artisan queue:work --queue=work2"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/queue-work2.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
priority=4
directory=/app
stopwaitsecs=10

[program:scout]
process_name=scout_%(program_name)s_%(process_num)02d
command=/usr/local/bin/run-worker-loop.sh "php /app/artisan queue:work --queue=scout --tries=3"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/queue-scout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
numprocs=3
process_name=%(program_name)s_%(process_num)02d
stopasgroup=true
killasgroup=true
priority=4
directory=/app
stopwaitsecs=10
